# backend/rfid-kiosk.service
[Unit]
Description=RFID-Controlled EV Charging Kiosk Backend Service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/ev_charging_kiosk/backend
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="PORT=5000"
ExecStart=/usr/bin/python3 /home/pi/ev_charging_kiosk/backend/rfid_service.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target


